# roles/setup/templates/suspendfix.service.tmpl
[Unit]
Description=Run user script after suspend

[Service]
ExecStart=/opt/sleep.sh

[Install]
WantedBy=multi-user.target
